# data/supervisor.conf
[group:sync_repo]
programs=sync_repo_1,sync_repo_2

[program:sync_repo_1]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_user_repo -e=prod banditore.sync_user_repo
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_user_repo_1.err
stdout_logfile=/space/logs/supervisord/sync_user_repo_1.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:sync_repo_2]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_user_repo -e=prod banditore.sync_user_repo
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_user_repo_2.err
stdout_logfile=/space/logs/supervisord/sync_user_repo_2.log
user=deploy
environment = http_proxy="",https_proxy=""
