# data/supervisor.conf
[group:sync_repo]
programs=sync_repo_1,sync_repo_2

[program:sync_repo_1]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_starred_repos -e prod banditore.sync_starred_repos
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_starred_repos_1.err
stdout_logfile=/space/logs/supervisord/sync_starred_repos_1.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:sync_repo_2]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_starred_repos -e prod banditore.sync_starred_repos
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_starred_repos_2.err
stdout_logfile=/space/logs/supervisord/sync_starred_repos_2.log
user=deploy
environment = http_proxy="",https_proxy=""

[group:sync_repo]
programs=sync_version_1,sync_version_2,sync_version_3,sync_version_4

[program:sync_version_1]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_versions -e prod banditore.sync_versions
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_versions_1.err
stdout_logfile=/space/logs/supervisord/sync_versions_1.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:sync_version_2]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_versions -e prod banditore.sync_versions
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_versions_2.err
stdout_logfile=/space/logs/supervisord/sync_versions_2.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:sync_version_3]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_versions -e prod banditore.sync_versions
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_versions_3.err
stdout_logfile=/space/logs/supervisord/sync_versions_3.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:sync_version_4]
directory=/path/to/banditore
command=/usr/bin/php bin/console swarrot:consume:banditore.sync_versions -e prod banditore.sync_versions
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/sync_versions_4.err
stdout_logfile=/space/logs/supervisord/sync_versions_4.log
user=deploy
environment = http_proxy="",https_proxy=""
